# axiom_beta_control_daemon/systemd/axiom.socket
[Unit]
Description=Axiom Socket
PartOf=axiom.service

[Socket]
ListenDatagram=/tmp/axiom_daemon

[Install]
WantedBy=sockets.target
